# files/systemd/tmux-control-plane.service
[Unit]
Description=tmux-control-plane MCP server
After=network.target claude-tmux.service

[Service]
Type=simple
ExecStart=%h/tmux-control-plane/start.sh
Restart=always
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=default.target
