# templates/github-mcp-watchdog.service.j2
[Unit]
Description=GitHub MCP health check (watchdog)

[Service]
Type=oneshot
ExecStart=%h/github-mcp/healthcheck.sh
EnvironmentFile=%h/.config/github-mcp/env
