# templates/github-mcp.service.j2
[Unit]
Description=GitHub MCP server (via mcp-proxy)
After=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
ExecStart=%h/github-mcp/start.sh
Restart=always
RestartSec=5
EnvironmentFile=%h/.config/github-mcp/env

[Install]
WantedBy=default.target
